# dfaprs/dfaprs-upstartd.conf
start on runlevel [2345]
stop on runlevel [!2345]

#expect fork

respawn
respawn limit 5 5

script
    /usr/local/bin/dfaprs | logger -t dfaprs    
end script

post-start script
end script

pre-start script
        logger -t dfaprs Starting dfaprs
end script

pre-stop script
        logger -t dfaprs Stopping dfaprs
end script
